# timeline-app/docker/supervisor/processes.conf
[supervisord]
nodaemon=true

[program:apache]
command=apache2-foreground
autostart=true
numprocs=1
